# mac-mini/Library/LaunchAgents/com.greg.copyparty-cloudflared.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.greg.copyparty-cloudflared</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/greg.linscheid/.local/bin/copyparty-cloudflared</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/greg.linscheid/Library/logs/copyparty/copyparty-cloudflared.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/greg.linscheid/Library/logs/copyparty/copyparty-cloudflared.error.log</string>
</dict>
</plist>
