# mac-mini/Library/LaunchAgents/com.greg.copyparty-cloudflared.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.greg.copyparty-mini</string>
    <key>ProgramArguments</key>
    <array>
        <string>zsh</string>
        <string>-c</string>
        <string>source ~/.zshenv/init/copyparty-env.sh; source ~/.zshenv/post/mac-mini.sh; copyparty-cloudflared</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/usr/local/var/log/copyparty-mini.log</string>
    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/copyparty-mini.error.log</string>
</dict>
</plist>
